# roles/mounts/files/mnt-linuxdata.automount
[Unit]
Description=Automount Optional "linuxgames" Disk Mount

[Automount]
Where=/mnt/linuxdata

[Install]
WantedBy=multi-user.target
