# roles/mounts/files/mnt-linuxgames.automount
[Unit]
Description=Automount Optional "linuxgames" Disk Mount

[Automount]
Where=/mnt/linuxgames

[Install]
WantedBy=multi-user.target
