# roles/mounts/files/mnt-windata.automount
[Unit]
Description=Automount Optional "windata" Disk Mount

[Automount]
Where=/mnt/windata

[Install]
WantedBy=multi-user.target
