# roles/mounts/files/mnt-Downloads.automount
[Unit]
Description=Automount HOMENAS Downloads Share

[Automount]
Where=/mnt/Downloads
TimeoutIdleSec=180

[Install]
WantedBy=multi-user.target
