# roles/mounts/files/mnt-data.automount
[Unit]
Description=Automount Optional "data" Disk Mount

[Automount]
Where=/mnt/data

[Install]
WantedBy=multi-user.target
