# files/dotfiles/Configs/systemd/niri-watcher.service
[Unit]
Description=Niri VRR/Fullscreen Watcher
Documentation=man:niri(1)
# Only start after niri.service is active
After=niri.service
# Bind this service's lifecycle to niri.service
BindsTo=niri.service
# Stop if niri.service stops (e.g., switching to another WM)
PartOf=niri.service
# Ensure we only run in niri sessions
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Type=simple
ExecStart=/usr/local/bin/niri-watcher
Restart=on-failure
RestartSec=5
KillMode=mixed
TimeoutStopSec=10
Environment=WAYLAND_DISPLAY=%E{WAYLAND_DISPLAY}

[Install]
# Only wanted when niri.service is wanted
WantedBy=niri.service
