# files/dotfiles/Configs/systemd/.config/systemd/user/vdirsyncer.service
[Unit]
Description=Sync calendars and contacts with vdirsyncer

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vdirsyncer sync
