# roles/base/files/systemd/beesd@.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/beesd --no-timestamps %i
Restart=on-failure
RestartSec=10s
SuccessExitStatus=124
IOSchedulingClass=idle
CPUSchedulingPolicy=idle
OOMScoreAdjust=500
