# systemd-automount/mnt-SSDDATA1.automount
[Unit]
Description=Automount SSDDATA1 disk

[Automount]
Where=/mnt/SSDDATA1

[Install]
WantedBy=multi-user.target
