# roles/base/files/systemd/beesd-maintenance.timer
[Unit]
Description=Run beesd deduplication during maintenance window

[Timer]
# Define your maintenance window start time (e.g., every day at 9:00 AM)
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
