# roles/base/files/systemd/beesd-maintenance.service
[Unit]
Description=Dynamically discover and launch all beesd configuration templates

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=60s
TimeoutStopSec=60s
ExecStart=/bin/bash -c '\
  for conf in /etc/bees/*.conf; do \
    uuid=$(basename "$conf" .conf); \
    echo "Starting beesd@$uuid..."; \
    systemctl start beesd@$uuid.service; \
  done'
ExecStop=/bin/bash -c '\
  for conf in /etc/bees/*.conf; do \
    uuid=$(basename "$conf" .conf); \
    systemctl stop beesd@$uuid.service; \
  done'

