# inventory/group_vars/local.yml
---
# System Configuration
timezone: America/Denver
graphics_card: amd

# Path Configuration
# Use inventory_dir | dirname to get the project root regardless of which playbook is run
# inventory_dir always points to the inventory/ directory
project_root: "{{ inventory_dir | dirname }}"
support_dir: "{{ project_root }}/files/support"
dotfiles_root: "{{ project_root }}/files/dotfiles"
nix_flake_root: "{{ home_realpath }}/.config/flakeroot"

# User Configuration
vault_file: "{{ support_dir }}/gpg_vault.yml"
home_dir: "{{ home_realpath }}"
scripts_dir: "{{ home_dir }}/.local/bin/scripts"
gnupg_key_dir: "{{ support_dir }}/.ssh/gnupg-keys"
ssh_key_dir: "{{ support_dir }}/.ssh"
pipewire_dir: "{{ home_dir }}/.config/pipewire"

# Kernel Arguments
kernel_args: "split_lock_detect=off amdgpu.ppfeaturemask=0xfffd7fff amdgpu.dcdebugmask=0x410"

# Program Checks
program_checks:
  hypr: "hyprland"
  MangoHud: "mangohud"
  OpenRGB: "openrgb"
  topgrade.d: "topgrade"

# Script Links ($HOME/.local/bin/scripts/ exposed globally via /usr/local/bin/)
script_links:
  - src: "fish.sh"
    dest: "/usr/local/bin/fish.sh"
  - src: "nscb.pyz"
    dest: "/usr/local/bin/nscb"
  - src: "protonfetcher.pyz"
    dest: "/usr/local/bin/protonfetcher"
  - src: "urh.pyz"
    dest: "/usr/local/bin/urh"
  - src: "squish.sh"
    dest: "/usr/local/bin/squish"
  - src: "unsquish.sh"
    dest: "/usr/local/bin/unsquish"
  - src: "khalreminder.sh"
    dest: "/usr/local/bin/khalreminder"
  - src: "niri_watcher.py"
    dest: "/usr/local/bin/niri-watcher"
  - src: "gamemode.pyz"
    dest:
      - "/usr/local/bin/gamemode"
      - "{{ home_dir }}/.local/bin/gamemode"

# Program Handlers (run after dotfiles are processed)
program_handlers:
  bat:
    - command: "bat cache --build"
      creates: "{{ home_dir }}/.cache/bat"

# Feature Flags
enable_lazyvim: true
enable_flatpak: true
remove_existing_flatpaks: true
install_firefox: false
install_peazip: true
install_arpc_fix: false
enable_fix_wakeups: false

# Systemd user services to enable (by unit name, not full path)
enabled_systemd_units: []
