# files/dotfiles/Configs/niri/.config/systemd/user/kde-kwallet-init.service
[Unit]
Description=Initialize KDE Wallet via PAM

[Service]
Type=oneshot
ExecStart=/usr/libexec/pam_kwallet_init
RemainAfterExit=yes
